# recipes-connectivity/kuksa-val/kuksa-dbc-feeder/kuksa-dbc-feeder.service
[Unit]
Description=Eclipse KUKSA.val DBC feeder
Requires=kuksa-databroker.service can-dev-helper.service
After=kuksa-databroker.service can-dev-helper.service

[Service]
ExecStart=/usr/bin/dbcfeeder.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
